# 0-学习环境的搭建/build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/devops/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/devops/supervisord.pid

[program:jupyterlab]
command=conda run -n ai-learning-notes jupyter-lab --allow-root --no-browser --NotebookApp.allow_remote_access=True --ip=0.0.0.0 --NotebookApp.token='UserPas$0w'
autostart=true
stderr_logfile=/devops/jupyterlab.err.log
stdout_logfile=/devops/jupyterlab.out.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
stderr_logfile=/devops/sshd.err.log
stdout_logfile=/devops/sshd.out.log